# traits/mods/vfd_lcd/files/etc/rc.local
/usr/bin/lcd-write "This is the Minikrebs   Version 84"
